# project-gamera-client.cron
#!/bin/sh
# project-gamera-client  A shell script to generate shared configuration data
# for a Project Gamera Server 

/usr/share/project-gamera/pg-client.sh
